# systemd/irig-sender.service
[Unit]
Description=IRIG-H Timecode Sender
After=multi-user.target

[Service]
Type=simple
# Default pins: BCM GPIO 11 (normal), inverted disabled
# To customize: ExecStart=/usr/local/bin/irig_sender -p 17 -n 27
ExecStart=/usr/local/bin/irig_sender
Restart=always
RestartSec=10
User=root
Nice=-20
StandardOutput=journal
StandardError=journal

# Create a working directory for output files
WorkingDirectory=/var/log/irig-sender
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Ensure the service stops gracefully
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target